# init-d-gentoo-funtoo-botcha
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License v2

USER="morfetico"

depend() {
need netif.eth0


}

start() {
ebegin "Iniciando o Botcha \"${USER}\""
su ${USER} -c ~${USER}/botcha/start.sh

}

stop() {
ebegin "Parando o Botcha \"${USER}\""
CHECK=$(pgrep -u ${USER} python2.6|wc -l)
if [ "$CHECK" -ne 0 ]; then
kill -9 $(pgrep -u ${USER} python2.6)
fi

}

status() {
ebegin "Processo em execução"
echo -e "USER: ${USER}\nPID: $(pgrep -u ${USER} python2.6)\n\n$(ps -uwp $(pgrep -u ${USER} python2.6))"
}
